# etc/services-config/hadoop-bootstrap.conf
[program:hadoop-bootstrap]
priority = 5
command = /usr/sbin/hadoop-bootstrap
autostart = true
user = hadoop
startsecs = 0
startretries = 0
autorestart = false
redirect_stderr = true
stdout_logfile = /var/log/hadoop-bootstrap.log
stdout_events_enabled = true
